# switcher2_bridge/rootfs/etc/services.d/switcher2_bridge/run
#!/command/with-contenv bashio
set -e

exec python3 /opt/switcher2_bridge/addon_main.py
